# modules/home/shell/.config/profile.d/editor.profile
export EDITOR="emacsclient -t"
export VISUAL="emacsclient -c -a ''"
export ALTERNATE_EDITOR=""
export GIT_EDITOR="emacsclient -t"
