# modules/home/shell/.config/profile.d/hub.profile
# useful in Emacs (eg auto-insert for .el)
export ORGANIZATION="Hubert Behaghel"
